# volumio/lib/systemd/system/fix-alpha-hosts.service
[Unit]
Description=Fix /etc/hosts and ensure test sentinels

[Service]
Type=oneshot
ExecStart=/bin/fix-alpha-hosts.sh
