# volumio/lib/systemd/system/fix-alpha-hosts.path
[Unit]
Description=Watch /etc/hosts for changes

[Path]
PathChanged=/etc/hosts

[Install]
WantedBy=multi-user.target
