# volumio/lib/systemd/system/ip-changed@.target
[Unit]
# Part of volumio-ipchange monitoring service
Description=IP Address changed on %i
