# volumio/lib/systemd/system/setdatetime-helper.timer.service
[Unit]
Description=Volumio Time Sync Watchdog Service

[Service]
Type=oneshot
ExecStart=/usr/bin/setdatetime-helper.sh
