# volumio/lib/systemd/system/volumio_fs_protection.service
# Systemd service to manage fs.protected_regular and fs.protected_fifos settings
# This service ensures that the kernel security settings are enforced or disabled as needed
# to maintain compatibility with applications that require access to world-writable directories.
#
# ExecStart disables the settings, allowing legacy applications that need shared access.
# ExecStop enables the settings, improving security by preventing privilege escalation attacks.

[Unit]
Description=Manage fs.protected_regular and fs.protected_fifos for Volumio
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/volumio_fs_protection.sh disable
ExecStop=/bin/volumio_fs_protection.sh enable
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
