# volumio/etc/systemd/system/nmbd.service.d/override.conf
[Unit]
After=winbind.service
Requires=winbind.service
