# volumio/etc/systemd/system/alsa-restore.service.d/override.conf
[Service]
ExecStartPre=/bin/mkdir -p /run/alsa /run/alsa/runtime
ExecStartPre=/usr/bin/wait-for-cards.sh
