# volumio/etc/systemd/system/winbind.service.d/override.conf
[Unit]
After=network-online.target smbd.service
Requires=network-online.target smbd.service
