# volumio/etc/systemd/system/smbd.service.d/override.conf
[Service]
ExecStartPre=/bin/mkdir -m 700 -p /var/log/samba/cores
